# roles/proxmox_node/templates/DEPRECATED_ups-battery-monitor.timer.j2
[Unit]
Description=Run UPS battery monitor every minute

[Timer]
# First run 15s after boot/activation, then every 60s thereafter
OnActiveSec=15s
OnBootSec=15s
OnUnitActiveSec=60s
AccuracySec=10s
Persistent=true
Unit=ups-battery-monitor.service

[Install]
WantedBy=timers.target
