# roles/proxmox_node/templates/radeontop-exporter.service.j2
[Unit]
Description=Export radeontop metrics for Prometheus

[Service]
Type=oneshot
ExecStart={{ radeontop_exporter_path }}
