# roles/share_drive_probe/templates/share_drive_probe.service.j2
# vim: ft=systemd

[Unit]
Description=Share Drive Probe (SMB/NFS liveness)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart={{ share_drive_probe_bin }}
TimeoutStartSec=90s

# Treat “1” as success (partial failure already lands in Prometheus)
SuccessExitStatus=0 1

User=root
Group=root
UMask=0022
# Nice = 5
Nice=1

StandardOutput=journal
StandardError=journal
SyslogIdentifier=share-drive-probe

Environment=SDP_TZ=Europe/Amsterdam
Environment=SDP_MODE=auto
Environment=SDP_RO_START={{ share_drive_probe_quiet_hours_start }}
Environment=SDP_RO_END={{ share_drive_probe_quiet_hours_end }}
Environment=SDP_LOG_LEVEL=info
