# roles/media_vm/templates/sleep_contains/docker_quiet@.service
[Unit]
Description=Quiet hours: %i Docker containers
Wants=docker.service
After=docker.service

[Service]
Type=oneshot
# Run as root unless your service account is in the docker group
ExecStart=/usr/local/bin/docker-quiet.sh %i
