# roles/media_vm/templates/sleep_contains/docker_quiet@unpause.timer.j2
[Unit]
Description=Unpause media containers at 10:00

[Timer]
OnCalendar=10:00
Persistent=true

[Install]
WantedBy=timers.target
