# roles/media_vm/templates/sleep_contains/docker_quiet@pause.timer
[Unit]
Description=Pause media containers at 22:00

[Timer]
OnCalendar=22:00
Persistent=true

[Install]
WantedBy=timers.target
