# roles/media_vm/templates/sleep_containers/docker-sleep@pause.timer.j2
# vim: ft=systemd

[Unit]
Description=Pause media containers at {{ docker_quiet_hours_start }} and retry periodically during quiet hours

[Timer]
OnCalendar={{ docker_quiet_hours_start }} Europe/Amsterdam
OnUnitActiveSec={{ docker_quiet_pause_retry_minutes | default(10) }}min
Persistent=false
Unit=docker-sleep@pause.service

[Install]
WantedBy=timers.target
