# roles/media_vm/templates/sleep_containers/docker-sleep@unpause.timer.j2
[Unit]
Description=Unpause media containers at {{ docker_quiet_hours_end }}

[Timer]
OnCalendar={{ docker_quiet_hours_end }} Europe/Amsterdam
Persistent=false

[Install]
WantedBy=timers.target
