# roles/media_vm/templates/sleep_containers/docker-quiet@pause.timer.j2
[Unit]
Description=Pause media containers at {{ docker_quiet_hours_start }}

[Timer]
OnCalendar={{ docker_quiet_hours_start }}
Persistent=true

[Install]
WantedBy=timers.target
