# roles/media_vm/templates/sleep_containers/docker-sleep@pause.timer.j2
[Unit]
Description=Pause media containers at {{ docker_quiet_hours_start }}

[Timer]
OnCalendar={{ docker_quiet_hours_start }} Europe/Amsterdam
Persistent=false

[Install]
WantedBy=timers.target
